# hw1/user/srm-empty.timer
[Unit]
Description=Run Smart Trash Can cleanup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

