# user/srm-empty.service
[Unit]
Description=Smart Trash Can cleanup service

[Service]
Type=oneshot
ExecStart=%h/bin/srm-empty

